# backup.service
[Unit]
Description=Start Backup-Script

[Service]
Type=oneshot
#Pfad anpassen
ExecStart=/opt/tar-Backup.sh
